# templates/docker-container.service.j2
[Unit]
Description={{ item.value.description }}
{% if shared_storage %}
After=docker.service autofs.service {{ data_mount_root | regex_replace('^\\/', '') }}.mount
Requires=docker.service autofs.service {{ data_mount_root | regex_replace('^\\/', '') }}.mount {{ item.value.depends_on | default([]) | join(' ') }}
RequiresMountsFor={{ data_mount_root }}
ConditionPathIsMountPoint={{ data_mount_root }}
{% else %}
After=docker.service
Requires=docker.service {{ item.value.depends_on | default([]) | join(' ') }}
{% endif %}

[Service]
Restart=always
ExecStart=/usr/bin/docker start -a {{ item.key }}
ExecStop=/usr/bin/docker stop -t 2 {{ item.key }}

[Install]
WantedBy=multi-user.target
